# directslave.service
[Unit]
After=directslave.service

[Service]
ExecStart=/usr/local/directslave/bin/directslave --run

[Install]
WantedBy=default.target
